# etc/profile.d/strapenv
# WARNING
#
# DO NOT MODIFY THIS FILE. IT MAY BE AUTOMATICALLY RE-GENERATED BY STRAP AT ANY TIME.
#
# If you want to modify, overwrite or unset anything in here, do that in your
# ~/.bash_profile or ~/.zshenv file *after* the line that sources $HOME/.strap/strapenv
#
